# cookbooks/tw-mvp/templates/default/mvp_frontend.service.erb
[Unit]
Description=Java Service of MVP Front-End

[Service]
User=mvp
EnvironmentFile=/etc/mvp_frontend_env
ExecStart=/usr/bin/java -jar /srv/mvp/build/front-end.jar
ExecStop=/bin/kill -TERM $MAINPID
TimeoutSec=300
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
