# cookbooks/tw-mvp/templates/default/mvp_newsfeed.service.erb
[Unit]
Description=Java Service of MVP Quote

[Service]
User=mvp
Environment=APP_PORT=<%= node['tw-mvp']['newsfeed_service_port'] %>
ExecStart=/usr/bin/java -jar /srv/mvp/build/newsfeed.jar
ExecStop=/bin/kill -TERM $MAINPID
TimeoutSec=300
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
